# davis-wx.service
[Unit]
Description=Davis Weather Station Application
After=network.target

[Service]
Type=simple
User=davis
Group=davis
WorkingDirectory=/opt/davis-wx/backend
ExecStart=/opt/davis-wx/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
Environment=KANFEI_SERIAL_PORT=/dev/ttyUSB0
Environment=KANFEI_BAUD_RATE=19200
Environment=KANFEI_DB_PATH=/var/lib/davis-wx/weather.db

[Install]
WantedBy=multi-user.target
